# monitoring.service
[Unit]
Description=Monitoring service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitoring.sh
